# systemd/oled-dashboard.service
[Unit]
Description=OLED Dashboard - Web-configurable OLED display manager
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/oled-dashboard
ExecStart=/opt/oled-dashboard/venv/bin/python -m oled_dashboard
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
